# deploy/ironroll.service
[Unit]
Description=IronRoll Dice Bot Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ironroll
ExecStart=/opt/ironroll/ironroll
Restart=always
RestartSec=5
EnvironmentFile=/opt/ironroll/.env

[Install]
WantedBy=multi-user.target
